# scripts/kt8900copilot-api.service
[Unit]
Description=KT8900 Copilot API Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kt8900copilot/server
ExecStart=/usr/bin/python3 /opt/kt8900copilot/server/api_server.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
